# flask-service.service
[Unit]
Description=Flask Unified Service
After=network.target
Wants=network.target

[Service]
Type=forking
User=root
Group=root
WorkingDirectory=/root/dify_tools
Environment=PATH=/root/miniconda3/envs/venv312/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=CONDA_DEFAULT_ENV=venv312
Environment=FLASK_ENV=production
Environment=FLASK_HOST=0.0.0.0
Environment=FLASK_PORT=5000
Environment=WORKERS=4
Environment=TIMEOUT=30
Environment=LOG_LEVEL=info
Environment=LOG_DIR=./logs
Environment=PID_FILE=/tmp/flask_service.pid

# 启动命令
ExecStart=/root/dify_tools/start_production.sh start
ExecStop=/root/dify_tools/start_production.sh stop
ExecReload=/root/dify_tools/start_production.sh restart

# 进程管理
PIDFile=/tmp/flask_service.pid
Restart=always
RestartSec=10
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# 安全设置
NoNewPrivileges=false
PrivateTmp=false

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=flask-service

[Install]
WantedBy=multi-user.target
